# overlays/system/etc/profile.d/prompt.sh
export PS1='\[\e[1;37m\]\u-\[\e[m\]$(pwd) '
